# sulad/files/etc/init.d/sulad
#!/bin/bash
### BEGIN INIT INFO
# Provides:          sulad
# Required-Start:    $local_fs $network $named $time $syslog
# Required-Stop:     $local_fs $network $named $time $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       <DESCRIPTION>
### END INIT INFO

CONFIG_FILE="/etc/sulad.conf"
PIDFILES="/var/run/sulad.*"
cmd="sulad"
GPIO=71

gpio_check(){
    if  [ ! -e /sys/class/gpio/gpio$GPIO/value ]; then
      echo $GPIO > /sys/class/gpio/export
      echo out > /sys/class/gpio/gpio$GPIO/direction
    fi
}

start() {
    echo 'Starting service…' >&2
    IFS=$'\n'
    for var in $(cat $CONFIG_FILE)
    do   
        echo "$var --pid-file /var/run/sulad" | xargs sulad
    done
    echo 'Service started' >&2
<<<<<<< HEAD:sulad/files/etc/init.d/sulad
    
    gpio_check
    
    sh -c "echo 1 > /sys/class/gpio/gpio$GPIO/value"
    
=======
    if [ -d /sys/class/gpio/gpio71 ]; then
      sh -c 'echo 1 > /sys/class/gpio/gpio71/value';
    fi
>>>>>>> cd263fd78db013bc262b0f9161e4addc39dc4f44:sulad/files/etc/init.d/sulad
}

stop() {
  echo $(cat $PIDFILES) | xargs kill
  echo 'Service stopped' >&2
<<<<<<< HEAD:sulad/files/etc/init.d/sulad
  
  gpio_check
    
  sh -c "echo 0 > /sys/class/gpio/gpio$GPIO/value"

=======
  if [ -d /sys/class/gpio/gpio71 ]; then
    sh -c 'echo 0 > /sys/class/gpio/gpio71/value';
  fi
>>>>>>> cd263fd78db013bc262b0f9161e4addc39dc4f44:sulad/files/etc/init.d/sulad
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 3
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
esac